# roles/hermes/templates/hermes-serve.service.j2
[Unit]
Description=Hermes Agent Dashboard Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart={{ hermes_venv_python | dirname }}/hermes serve --host 0.0.0.0 --port {{ hermes_dashboard_port }}
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
